# config/panelkit.service.template
[Unit]
Description=PanelKit UI Application
After=network.target

[Service]
Type=simple
User=__USER__
WorkingDirectory=__WORKDIR__
ExecStart=/bin/bash -c "PANELKIT_USE_RENDERING=1 RUST_LOG=info __WORKDIR__/panelkit"
Restart=on-failure
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target