# panelkit_sdl/deploy/panelkit.service
[Unit]
Description=PanelKit Touch Panel Application
After=network.target graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
ExecStartPre=/bin/mkdir -p /var/log/panelkit
ExecStart=/tmp/panelkit/panelkit
Restart=always
RestartSec=5
User=root
StandardOutput=append:/var/log/panelkit/panelkit.log
StandardError=append:/var/log/panelkit/panelkit.log
Environment=SDL_FBDEV=/dev/fb0
Environment=SDL_VIDEODRIVER=fbcon
Environment=SDL_AUDIODRIVER=dummy
TimeoutStartSec=30
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target